# feeds/lorawan-gateway/luci-m2-allinone-general/root/etc/uci-defaults/31_luci-theme-bootstrap
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci get luci.themes.SenseCAP >/dev/null 2>&1 || \
	uci batch <<-EOF
		set luci.themes.SenseCAP=/luci-static/bootstrap
		set luci.main.mediaurlbase=/luci-static/bootstrap
		commit luci
	EOF
fi

exit 0
