# feeds/sensecap/ubus-serve/files/ubus-serve.init
#!/bin/sh /etc/rc.common
START=51
STOP=10

USE_PROCD=1
PROG=/usr/bin/ubus-serve

start_service()
{
	procd_open_instance
	procd_set_param command $PROG
	procd_set_param respawn 3600 60 60
	procd_close_instance
}

stop_service() {
	logger "ubus-serve: Stop"
}

reload_service()
{
	stop
	start
}
